# anbox-session-manager.service
[Unit]
Description=Anbox Session Manager

[Service]
ExecStart=/usr/bin/anbox session-manager --gles-driver=host

[Install]
WantedBy=default.target
